# configs/puppet/modules/automaton/templates/shmwatch.cron
#!/bin/sh

flags=-umc

/usr/sbin/tmpwatch "$flags" -d -f 2 <%= @dir_www %>/tmp
