# configs/puppet/modules/zloy/files/nodejs.conf
description "node.js server"

start on filesystem

task

exec /nodejs.sh
